# supervisord.d/0_stand.ini
[program:example]
command=/opt/src/elasticsearch-2.4.6/bin/elasticsearch
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/opt/src/elasticsearch-2.4.6/bin/elasticsearch/elastic.out
stderr_logfile=/opt/src/elasticsearch-2.4.6/bin/elasticsearch/elastic.err
