# supervisord.d/1_kafka_zookeeper.ini
[program:kafka]
command=/opt/src/kafka_2.11-0.11.0.0/bin/kafka-server-start.sh /opt/src/kafka_2.11-0.11.0.0/config/server.properties
autostart=true
autorestart=unexpected
startsecs=20
startretries=3
stdout_logfile=/mnt/data/kafka_2.11-0.11.0.0/kafka.out
stderr_logfile=/mnt/data/kafka_2.11-0.11.0.0/kafka.err

[program:zookeeper]
command=/opt/src/zookeeper-3.4.10/bin/zkServer.sh start-foreground
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/zookeeper-3.4.10/zk.out
stderr_logfile=/mnt/data/zookeeper-3.4.10/zk.err
