# supervisord.d/2_logstash.ini
[program:logstash_errorhandling]
command=/opt/src/logstash-5.5.2/bin/logstash -f /opt/src/logstash-5.5.2/conf.d/error_handling.conf --path.data=/mnt/data/logstash-5.5.2/data/errorhandling
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/logstash-5.5.2/logs/logstash_errorhandling.out
stderr_logfile=/mnt/data/logstash-5.5.2/logs/logstash_errorhandling.err

[program:logstash_drmtoday]
command=/opt/src/logstash-5.5.2/bin/logstash -f /opt/src/logstash-5.5.2/conf.d/drm_today.conf --path.data=/mnt/data/logstash-5.5.2/data/drmtoday
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/logstash-5.5.2/logs/logstash_drmtoday.out
stderr_logfile=/mnt/data/logstash-5.5.2/logs/logstash_drmtoday.err

[program:logstash_stream]
command=/opt/src/logstash-5.5.2/bin/logstash -f /opt/src/logstash-5.5.2/conf.d/stream.conf --path.data=/mnt/data/logstash-5.5.2/data/stream
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/logstash-5.5.2/logs/logstash_stream.out
stderr_logfile=/mnt/data/logstash-5.5.2/logs/logstash_stream.err

[program:logstash_mbf]
command=/opt/src/logstash-5.5.2/bin/logstash -f /opt/src/logstash-5.5.2/conf.d/mbf.conf --path.data=/mnt/data/logstash-5.5.2/data/mbf
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/logstash-5.5.2/logs/logstash_mbf.out
stderr_logfile=/mnt/data/logstash-5.5.2/logs/logstash_mbf.err

[program:logstash_vmx]
command=/opt/src/logstash-5.5.2/bin/logstash -f /opt/src/logstash-5.5.2/conf.d/vmx.conf --path.data=/mnt/data/logstash-5.5.2/data/vmx
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/logstash-5.5.2/logs/logstash_vmx.out
stderr_logfile=/mnt/data/logstash-5.5.2/logs/logstash_vmx.err
