# supervisord.d/3_flume.ini
[program:flume_useractivity]
command=/opt/src/apache-flume-1.7.0-bin/bin/flume-ng agent -n FKafka -c conf -f /opt/src/apache-flume-1.7.0-bin/conf/fk.properties.template
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stdout_logfile=/mnt/data/apache-flume-1.7.0-bin/flume_useractivity.out
stderr_logfile=/mnt/data/apache-flume-1.7.0-bin/flume_useractivity.err
